# tools/com.beako.utils.directhwd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.beako.utils.directhwd</string>

	<!-- Loader script (installed by `make install-daemon`). -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/sbin/directhwd</string>
	</array>

	<!-- Run once at boot; loading a kext needs root, so this is a LaunchDaemon
	     (system domain), not a per-user LaunchAgent. -->
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>

	<key>StandardOutPath</key>
	<string>/var/log/directhwd.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/directhwd.log</string>
</dict>
</plist>
